# cpu-usage-led.service
[Service]
ExecStart=python3 /usr/bin/cpu-usage-led.py

[Install]
WantedBy=default.target
