# cpu-usage-led.service
[Service]
ExecStart=/usr/local/bin/cpu-usage-led.py

[Install]
WantedBy=default.target
